# portfw.service
[Unit]
Description=Port Forward Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/portfw/portfw.env
ExecStart=/bin/bash -c 'source /etc/portfw/portfw.env 2>/dev/null || true; /usr/local/bin/portfw-restore.sh'

[Install]
WantedBy=multi-user.target
